# xdg-config/autostart/openvpn3-applet.desktop
[Desktop Entry]
Exec=sh -c "~/.openvpn3-applet/openvpn3applet.sh"
Icon=
Name=openvpn3-applet
Terminal=false
Type=Application
